# MAKE.list
/home/kahsolt/.byobu/
/home/kahsolt/.config/filezilla/
/home/kahsolt/.config/i3/
/home/kahsolt/.config/i3status/
/home/kahsolt/.config/pcmanfm/
/home/kahsolt/.config/sakura/
/home/kahsolt/.dbshell
/home/kahsolt/.elinks/
/home/kahsolt/.gitconfig
/home/kahsolt/.git-credentials
/home/kahsolt/.hosts
/home/kahsolt/.irssi/
/home/kahsolt/.moc/
/home/kahsolt/.mongorc.js
/home/kahsolt/.oh-my-zsh/
/home/kahsolt/.pip/
/home/kahsolt/.Rprofile
/home/kahsolt/.ssh/
/home/kahsolt/.tmux.conf
/home/kahsolt/.vnc/
/home/kahsolt/.wget-hsts
/home/kahsolt/.zlog*
/home/kahsolt/.zsh*
/home/kahsolt/.*rc
/home/kahsolt/.*_history
/etc/hostname
/etc/hosts
/etc/fstab
/etc/init.d/sslocal
/etc/init.d/ssserver
/etc/init.d/tightvncserver
/etc/apt/sources.list
/etc/apt/sources.list.d/
/etc/polipo/config
/etc/shadowsocks/config.json
/etc/mongodb.conf
/etc/mysql/conf.d/mysql.cnf